# turbo-switcher.service.example
[Unit]
Description=Enable turbo mode based on AC availability

[Service]
Type=simple
ExecStart=/path/to/turbo-switcher.sh
User=root
