# turbo-switcher.service.example
[Unit]
Description=Enable turbo mode based on AC availability

[Service]
Type=simple
ExecStart=/home/user/.config/systemd/user/scripts/turbo-switcher.sh
User=root
